# systemd/coopdoor-apply-schedule.service
[Unit]
Description=Apply Coop Door schedule (recompute timers)
After=coopdoor-api.service

[Service]
Type=oneshot
ExecStart=/opt/coopdoor/.venv/bin/python3 /opt/coopdoor/schedule_apply.py
